# server/deploy/yardguard-email.service
[Unit]
Description=YardGuard quote-form email endpoint (Resend)
After=network.target

[Service]
Type=simple
# --- Adjust these paths to where you cloned the repo on the Droplet ---
WorkingDirectory=/var/www/yardguard/server
# Run `which node` on the Droplet and use that path (may be /usr/bin/node,
# /usr/local/bin/node, or an nvm path). Node 18+ required for --env-file.
ExecStart=/usr/bin/node --env-file=/var/www/yardguard/server/.env /var/www/yardguard/server/index.js
Restart=on-failure
RestartSec=5
User=www-data
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
